# systemd/solar-system-alarm-deploy.service
[Unit]
Description=Deploy solar-system-alarm from GitHub main
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
User=opc
WorkingDirectory=/home/opc
Environment=APP_DIR=/opt/solar-system-alarm
Environment=REPO_URL=https://github.com/saqibroy/solar-system-alarm.git
Environment=BRANCH=main
Environment=COMPOSE_PROJECT_NAME=solar-system-alarm
ExecStart=/usr/local/bin/solar-system-alarm-deploy
